# es1/etc/udev/rules.d/99-monitor-hotplug.rules
KERNEL=="card0", ACTION=="change", SUBSYSTEM=="drm", ENV{HOTPLUG}=="1", ENV{XAUTHORITY}="/home/user/.Xauthority", RUN+="/home/user/bin/xrandr_LVDS"
